# n8n-evolution-api-odoo-18/crontab/crontab
0 2 * * * /backup/backup.sh >> /backup/backup.log 2>&1
